# target/add-ons/fs/etc/profile.d/bashrc.sh
export PS1='\h\$ '
#chrt -p -f 99 $$
#stty -ixon -ixoff -ixany < /dev/ttyO0
